# release/src/router/netatalk-3.0.5/distrib/initscripts/rc.netbsd
#!/bin/sh
#
# PROVIDE: netatalk
# KEYWORD: shutdown
#
# AFP fileserver for Macintosh clients (netatalk 3.0.5)
#

. /etc/rc.subr

name="netatalk"
rcvar=$name
command="/usr/sbin/netatalk"
etcdir="/usr/etc"
pidfile="/var/run/${name}.pid"
required_files="$etcdir/afp.conf"

load_rc_config $name
run_rc_command "$1"
